# deployment/rpi_ecu/systemd/autosar-ntp-time-provider.service
[Unit]
Description=Adaptive AUTOSAR NTP time synchronization provider
After=network-online.target
Wants=network-online.target
Before=autosar-platform-app.service autosar-exec-manager.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/autosar-ntp-time-provider
ExecStart=/usr/local/bin/autosar-ntp-time-provider-wrapper.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
