# deployment/rpi_ecu/systemd/autosar-platform-app.service
[Unit]
Description=Adaptive AUTOSAR platform process stack (EM/SM/PHM/Diag/Vehicle)
After=network-online.target autosar-iox-roudi.service autosar-vsomeip-routing.service autosar-time-sync.service autosar-persistency-guard.service autosar-iam-policy.service
Wants=network-online.target autosar-iox-roudi.service autosar-vsomeip-routing.service autosar-time-sync.service autosar-persistency-guard.service autosar-iam-policy.service
Before=autosar-exec-manager.service autosar-watchdog.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/autosar-platform-app
ExecStart=/usr/local/bin/autosar-platform-app-wrapper.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
