# deployment/rpi_ecu/systemd/autosar-ecu-full-stack.service
[Unit]
Description=Adaptive AUTOSAR ECU full-stack template (CAN + SOME/IP -> DDS)
After=network-online.target autosar-iox-roudi.service autosar-vsomeip-routing.service autosar-platform-app.service
Requires=autosar-iox-roudi.service autosar-platform-app.service
Wants=network-online.target autosar-iox-roudi.service autosar-vsomeip-routing.service autosar-platform-app.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/autosar-ecu-full-stack
ExecStart=/usr/local/bin/autosar-ecu-full-stack-wrapper.sh
Restart=always
RestartSec=2
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
