# deployment/rpi_ecu/systemd/autosar-sm-state.service
[Unit]
Description=Adaptive AUTOSAR state management daemon
After=network-online.target
Wants=network-online.target
Before=autosar-platform-app.service autosar-exec-manager.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/autosar-sm-state
ExecStart=/usr/local/bin/autosar-sm-state-wrapper.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
